# supervisord.conf
; supervisor config file

[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/var/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
childlogdir=/var/log/supervisor            ; ('AUTO' child log dir, default $TEMP)
; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///var/run//supervisor.sock ; use a unix:// URL  for a unix socket
; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.
[include]
;单个程序
files = /etc/supervisor/conf.d/chatbot.conf
;多个程序，同时写入conf.d文件夹中
;如有chatbot和face_recognition两个程序
;程序设计如下
;chatbot.conf和face_recognition.conf
files = /etc/supervisor/conf.d/*.conf
;注意，python中的文件后缀，可以更改
;经验：A用py2，B用py3
;使用的conf不兼容
;有后缀隔离自己配置：如*.ini, *.txt
;supervisor控制台
[inet_http_server]
# 监听的IP，形式为host：port
# 若port = *：8080报错，则使用0.0.0.0:9080
port=127.0.0.1:9090
;username=user
;password=123

[program:chatbot]
#执行gunicorn命令
command=/Users/xinxi/Documents/sndd/dockerflask/py2env/bin/gunicorn -w 1 -b 127.0.0.1:9999 chatbot:app
directory=/Users/xinxi/Documents/sndd/dockerflask
#必须添加用户
user=xinxi
#自动启动
autostart=true
#自动重启
autorestart=true
#日志
stdout_logfile=/Users/xinxi/Documents/sndd/dockerflask/logs/gunicorn_supervisor.log


[program:jenkins]
#执行gunicorn命令
command=java -jar jenkins.war
directory=/usr/local/Cellar/jenkins/2.186/libexec/
#必须添加用户
user=xinxi
#自动启动
autostart=true
#自动重启
autorestart=true
#日志
stdout_logfile=/Users/xinxi/Documents/sndd/dockerflask/logs/gunicorn_supervisor.log